# stx/toCOPY/aptly/supervisord.nginx.conf
[program:nginx]
command=/usr/sbin/nginx
redirect_stderr=true
stdout_logfile_maxbytes=1MB
