# stx/debian/common/toCOPY/aptly/supervisord.aptly.conf
[program:aptly]
command=/usr/bin/aptly api serve --no-lock --listen=":8080"
redirect_stderr=true
stdout_logfile_maxbytes=1MB
